# .config/systemd/user/newsboat.service
[Unit]
Description=Update newsboat feeds

[Service]
ExecStart=/usr/bin/env newsboat -x reload
